# build/bind/files/bind.xml
<?xml version="1.0"?>
<!DOCTYPE service_bundle SYSTEM "/usr/share/lib/xml/dtd/service_bundle.dtd.1">

<service_bundle type='manifest' name='ISCbind:named'>

<service
    name='network/dns/server'
    type='service'
    version='1'>

    <create_default_instance enabled='false' />

    <dependency
        name='filesystem_local'
        grouping='require_all'
        restart_on='none'
        type='service'>
            <service_fmri value='svc:/system/filesystem/local:default' />
    </dependency>

    <dependency
        name='multi_user'
        grouping='require_all'
        restart_on='none'
        type='service'>
            <service_fmri value='svc:/milestone/multi-user:default' />
    </dependency>

    <dependency
        name='loopback'
        grouping='require_any'
        restart_on='error'
        type='service'>
            <service_fmri value='svc:/network/loopback' />
    </dependency>

    <dependency
        name='network'
        grouping='optional_all'
        restart_on='error'
        type='service'>
            <service_fmri value='svc:/milestone/network' />
    </dependency>

    <dependent name='bind_multi-user-server'
        grouping='optional_all'
        restart_on='none'>
            <service_fmri value='svc:/milestone/multi-user-server' />
    </dependent>

    <exec_method
        type='method'
        name='start'
        exec='/opt/ISCbind/sbin/named -c %{config/file}'
        timeout_seconds='60'>
            <method_context security_flags='aslr'>
                <method_credential user='named' group='named'
                    privileges='basic,net_privaddr,!proc_info,!proc_session' />
            </method_context>
    </exec_method>

    <exec_method
        type='method'
        name='stop'
        exec=':kill'
        timeout_seconds='60' />

    <exec_method
        type='method'
        name='refresh'
        exec='/opt/ISCbind/sbin/rndc reload'
        timeout_seconds='60'>
            <method_context>
	        <method_credential user='named' group='named' />
            </method_context>
    </exec_method>

    <property_group name='config' type='application'>
        <propval name='file' type='astring'
            value='/etc/opt/ISCbind/named.conf' />
    </property_group>

    <stability value='Unstable' />

    <template>
        <common_name>
            <loctext xml:lang='C'>
                BIND DNS server
            </loctext>
        </common_name>
        <documentation>
            <manpage title='named' section='1M'
                manpath='/usr/man' />
        </documentation>
    </template>

</service>

</service_bundle>
